# C:/Users/Administrator/WorkBuddy/Claw/lvyou_smsweb_user/lvyou_smsweb/systemd/board-manager-v6.service
[Unit]
Description=Board Manager (IPv6)
After=network.target

[Service]
Type=simple
WorkingDirectory={{APPDIR}}
Environment=BMDB={{APPDIR}}/data/data.db
Environment=BMSTATIC={{APPDIR}}/static
Environment=BMDEVUSER={{SCANUSER}}
Environment=BMDEVPASS={{SCANPASS}}
Environment=BMUIUSER=admin
Environment=BMUIPASS={{UIPASS}}
ExecStart={{APPDIR}}/venv/bin/python -m uvicorn main:app --app-dir {{APPDIR}}/app --host :: --port {{APIPORT}}
Restart=always
RestartSec=2
User=root
Group=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
